# supervisor/supervisord.conf
[program:bcg]
command=bcg --device %(ENV_DEVICE)s
numprocs=1

[program:node-red]
command=node-red /tmp/node-red/.node-red/settings.js
numprocs=1